# supervisord/conf.d/hubot.conf
[program:hubot]
directory=/data/app
command=/data/app/bin/hubot-env
autorestart=true
stdout_events_enabled = true
stderr_events_enabled = true
